# roles/calibre-web/templates/cps.service.j2
Description=Calibre-{{ calibre_name }}

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python /var/www/calibre-{{ calibre_name }}/cps.py

[Install]
WantedBy=multi-user.target
